# deploy/systemd/epg-reconstruccion-nocturna.service
[Unit]
Description=EPG-UAC reconstruccion nocturna por identidad academica
After=epg-post-drive-identidades.service mariadb.service
[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/sistema_posgrado/backend
EnvironmentFile=-/opt/sistema_posgrado/backend/.env
ExecStart=/opt/sistema_posgrado/backend/venv/bin/python /opt/sistema_posgrado/backend/ejecutar_reconstruccion_nocturna.py
CPUQuota=90%
Nice=0
TimeoutStartSec=infinity
Restart=on-failure
RestartSec=15min
StartLimitIntervalSec=8h
StartLimitBurst=4
